# geek1011/root/etc/udev/rules.d/99-geek1011.rules
KERNEL=="loop0", RUN+="/usr/local/geek1011/init"
